# init/admin-drivers.conf
#!upstart
description "ninjablocks admin drivers"
author      "http://www.ninjablocks.com"

start on (local-filesystems and (net-device-up IFACE=eth0 or net-device-up IFACE=wlan0))
stop on shutdown

respawn
respawn limit 99 1

script
	export HOME="/opt/utilities"

	echo $$ > /var/run/admin-drivers.pid
	exec /opt/utilities/admin-drivers/bin/server --pid /var/run/ninjablock.pid -d /opt/ninja/drivers >> /var/log/admin-drivers.log 2>&1
end script

pre-start script
	echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (system) Starting." >> /var/log/admin-drivers.log
end script

pre-stop script
	rm /var/run/admin-drivers.pid
	echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (system) Stopping." >> /var/log/admin-drivers.log
end script
