# init/wifisetup.conf
description "ninjablocks wifi setup"
author      "http://www.ninjablocks.com"
 
# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on net-device-up
stop on shutdown
 
# Automatically Respawn:
respawn
respawn limit 99 5
 
script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/home/ubuntu"
 
    exec /opt/utilities/wifi/app.rb -p 80
end script